# src/filesystem/etc/cron.weekly/chkfxapt
#!/bin/sh

dash /usr/share/apps/q4os_system/bin/qapt_fix.sh --noninteractive 2>&1
exit 0
